# automations/adapters/launchd/template.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{label}}</string>
    <key>UserName</key>
    <string>agent-lab</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{repo_root}}/automations/bin/run-automation</string>
        <string>{{name}}</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{{repo_root}}</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/agent-lab</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
{{schedule}}
    <key>StandardOutPath</key>
    <string>{{repo_root}}/logs/{{log_name}}.log</string>
    <key>StandardErrorPath</key>
    <string>{{repo_root}}/logs/{{log_name}}.log</string>
    <key>RunAtLoad</key>
    <{{run_at_load}}/>
</dict>
</plist>
