# scripts/macos/com.tsd-agent-lab.stale-docs-check-full.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tsd-agent-lab.stale-docs-check-full</string>
    <key>UserName</key>
    <string>agent-lab</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/agent-lab/tsd-agent-lab/automations/bin/run-automation</string>
        <string>stale-docs-check-full</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/agent-lab/tsd-agent-lab</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/agent-lab</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>20</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/agent-lab/tsd-agent-lab/logs/stale-docs-check-full.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/agent-lab/tsd-agent-lab/logs/stale-docs-check-full.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
