# scripts/macos/com.tsd-agent-lab.health-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tsd-agent-lab.health-report</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/agent-lab/tsd-agent-lab/scripts/macos/health-report.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/agent-lab/tsd-agent-lab/logs/health-report.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/agent-lab/tsd-agent-lab/logs/health-report.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
